# archive/tqt3/etc/profile.d/tqt3.sh

# Set QTDIR if unset

test -n "${QTDIR+$QTDIR}" || QTDIR=/opt/trinity/lib/tqt3

export QTDIR

