# archive/logrotate/logrotate.cron
#! /bin/sh -
#
# /etc/cron.daily/logrotate - Rotates system logs periodically.
#

/usr/sbin/logrotate /etc/logrotate.conf
status=$?

if test $status -ne 0
then
    /usr/bin/logger \
     -t logrotate "ALERT exited abnormally with exit status $status"
fi
exit $status

