# archive/mlocate/etc/cron.daily/mlocate
#! /bin/sh -
#
# /etc/cron.daily/mlocate - Update the mlocate(1) database.
#

set -e
umask 022

ionice -c 3 nice -n 19 updatedb

