# archive/etc/rc.d/rc.local
#! /bin/sh -
#
# /etc/rc.d/rc.local - The local instructions for the boot process.
#

