# archive/man-db/etc/cron.daily/man-db
#! /bin/sh -
#
# /etc/cron.daily/man-db - Update manual page(s) database periodically.
#

set -e
umask 022

test -d /var/cache/man || mkdir -p /var/cache/man
ionice -c 3 nice -n 14 mandb --quiet

